# server/systemd/cyberalertx-generate.timer
# /etc/systemd/system/cyberalertx-generate.timer
#
# Fires cyberalertx-generate.service once an hour, at minute :00. Paired with
# `--limit 1` in the service, this mints at most one fresh story per hour
# (<=24/day) — a steady drip that keeps the feed current without flooding the
# Telegram channels (publish runs at :15, capped to one per channel).
#
# Persistent=true: if the VPS was off when a scheduled fire was due,
# systemd runs it as soon as the box is back up.
#
# RandomizedDelaySec=180: spread fires by up to 3 min so timers across
# the fleet don't all hit Anthropic at the same wall-clock second.

[Unit]
Description=CyberAlertX AI generate timer (hourly)
Requires=cyberalertx-generate.service

[Timer]
OnCalendar=*-*-* *:00:00
Persistent=true
RandomizedDelaySec=180
Unit=cyberalertx-generate.service

[Install]
WantedBy=timers.target
