# server/systemd/cyberalertx-generate.service
# /etc/systemd/system/cyberalertx-generate.service
#
# AI render — one-shot helper invoked by cyberalertx-generate.timer.
#
# Runs `generate --limit 2 --use-llm`. Idempotent: if the top-2 newest
# items are already in the AI cache, this is a no-op (cache-hits skip,
# zero API calls). Real cost is bounded by *new* items arriving since
# the last fire — i.e., by the actual news-cycle volume, not by how
# often the timer fires.
#
# The homepage feed cap (20 visible posts) is enforced at the API/feed
# layer, NOT here — old already-rendered posts accumulate in the cache
# and remain on the feed as new ones arrive. This service only mints
# 2 new renders per fire to keep API cost low.
#
# Why a separate service from cyberalertx-run:
#   * `run` is deterministic-only (no Anthropic ever). Keeping the AI
#     fire path in its own unit preserves that "AI is opt-in" invariant
#     at the systemd layer — easy to stop AI without stopping ingest.
#   * Failures here (rate limit, key expired, network) don't ripple
#     into the ingest service.
#
# Cost note: at 2 items × ~1.5 locales × ~$0.009/Haiku call, a fire
# costs <$0.03 in the busy case. Test cadence (every 6h = 4 fires/day):
# bounded by new-items-per-day × ~$0.013/item. Steady-state expectation
# ~$1–3/month. Easy to scale up to limit=3 every 4h once we see real
# numbers.

[Unit]
Description=CyberAlertX AI generate (one-shot, fired by timer)
After=network.target cyberalertx-api.service
Wants=cyberalertx-api.service

[Service]
Type=oneshot
User=cax
Group=cax
WorkingDirectory=/home/cax/cax
EnvironmentFile=/home/cax/cax/.env
ExecStart=/home/cax/cax/venv/bin/python -m cyberalertx.main generate --limit 4 --use-llm
# 4 items/fire (×4 fires/day = up to 16 renders/day) keeps the feed current
# with the news cycle. Renders go through the OAuth subscription (no API
# cost), so a higher limit isn't a billing concern — only a wall-clock one.
# At ~20-40s/locale, 4 items × ~1.5 locales ≈ 2-4 min; the 15-min cap leaves
# headroom for subscription latency before the timer gives up and retries.
TimeoutStartSec=900
StandardOutput=journal
StandardError=journal

# No [Install] — this is started by the .timer unit, not on boot.
