# server/systemd/cyberalertx-telegram.timer
# /etc/systemd/system/cyberalertx-telegram.timer
#
# Fires cyberalertx-telegram.service every 6 hours, 15 minutes AFTER the
# generate timer (which runs at 00,06,12,18:00). The 15-min offset gives
# cyberalertx-generate time to finish rendering the new batch, so the freshly
# rendered posts are publishable in the same cycle they were rendered.
#
# Persistent=true: if the VPS was off when a fire was due, systemd runs it as
# soon as the box is back up (the ledger prevents double-posting).
#
# RandomizedDelaySec=120: small jitter so fleet timers don't all hit the
# Telegram API at the same wall-clock second.

[Unit]
Description=CyberAlertX Telegram publish timer (every 6h, +15m after generate)
Requires=cyberalertx-telegram.service

[Timer]
OnCalendar=*-*-* 00,06,12,18:15:00
Persistent=true
RandomizedDelaySec=120
Unit=cyberalertx-telegram.service

[Install]
WantedBy=timers.target
