# server/systemd/cyberalertx-generate.service
# /etc/systemd/system/cyberalertx-generate.service
#
# AI render — one-shot helper invoked by cyberalertx-generate.timer.
#
# Runs `generate --limit 20 --use-llm`. Idempotent: cache-hits skip with
# zero API calls, so the limit is an *upper bound on new renders per
# fire*, not a fixed batch size. With store_cap=20 (cyberalertx/config.py:
# Settings.max_items_retained), matching the limit to the cap guarantees
# the entire store is AI-rendered — the homepage feed never has gaps
# because some stored items lack a render.
#
# Why a separate service from cyberalertx-run:
#   * `run` is deterministic-only (no Anthropic ever). Keeping the AI
#     fire path in its own unit preserves that "AI is opt-in" invariant
#     at the systemd layer — easy to stop AI without stopping ingest.
#   * Failures here (rate limit, key expired, network) don't ripple
#     into the ingest service.
#
# Cost note: per item ≈ 1.5 locales × ~$0.009/Haiku call ≈ $0.013. A
# fire's real cost = (new items since last fire) × $0.013, NOT 20×$0.013
# (cache-hits are free). Typical fire on a quiet news day: 0-3 new items
# → <$0.05. Worst case (store rotates fully in 6h, unlikely): ~$0.27.
# Steady-state expectation ~$2–5/month at 4 fires/day.

[Unit]
Description=CyberAlertX AI generate (one-shot, fired by timer)
After=network.target cyberalertx-api.service
Wants=cyberalertx-api.service

[Service]
Type=oneshot
User=cax
Group=cax
WorkingDirectory=/home/cax/cax
EnvironmentFile=/home/cax/cax/.env
ExecStart=/home/cax/cax/venv/bin/python -m cyberalertx.main generate --limit 20 --use-llm
# Each item renders in ~30-60s (1-2 locales × Haiku). A full 20-item
# backfill (only happens on first fire after a wipe / limit bump) is
# bounded by ~20 min; steady-state fires render 0-3 new items in <2 min.
# Beyond 20 min something's wrong (rate-limited, hung connection); the
# timer will retry on the next fire anyway.
TimeoutStartSec=1500
StandardOutput=journal
StandardError=journal

# No [Install] — this is started by the .timer unit, not on boot.
