# server/systemd/cyberalertx-telegram.timer
# /etc/systemd/system/cyberalertx-telegram.timer
#
# Fires cyberalertx-telegram.service once an hour at minute :15, 15 minutes
# AFTER the generate timer (which runs at :00). The offset gives
# cyberalertx-generate time to finish rendering the hour's story, so it's
# publishable in the same cycle. With `--limit 1` in the service this sends at
# most one message per channel per hour — pacing the channels, not flooding.
#
# Persistent=true: if the VPS was off when a fire was due, systemd runs it as
# soon as the box is back up (the ledger prevents double-posting).
#
# RandomizedDelaySec=120: small jitter so fleet timers don't all hit the
# Telegram API at the same wall-clock second.

[Unit]
Description=CyberAlertX Telegram publish timer (hourly, +15m after generate)
Requires=cyberalertx-telegram.service

[Timer]
OnCalendar=*-*-* *:15:00
Persistent=true
RandomizedDelaySec=120
Unit=cyberalertx-telegram.service

[Install]
WantedBy=timers.target
